# cellframe-node-asan.service
[Unit]
Description=Cellframe Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cellframe-node
ExecStartPre=+/bin/chown -R cellframe:cellframe /opt/cellframe-node
ExecStart=/opt/cellframe-node/bin/cellframe-node
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=always
User=cellframe
Group=cellframe
RestartSec=10
LogNamespace=cellframe
StandardOutput=null
StandardError=journal
Environment="ASAN_OPTIONS=strict_string_checks=1:detect_stack_use_after_return=1:detect_invalid_pointer_pairs=2:debug=1:atexit=1:abort_on_error=1:log_path=asanlog.txt"

# Hardening
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target